# mchrbl-bot.service
[Unit]
Description=MCHRBL Bot - Xiaomi Bootloader Unlock War Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/.openclaw/workspace/projects/mchrbl-bot
EnvironmentFile=/root/.openclaw/workspace/projects/mchrbl-bot/.env
ExecStart=/usr/bin/python3 /root/.openclaw/workspace/projects/mchrbl-bot/main.py
WorkingDirectory=/root/.openclaw/workspace/projects/mchrbl-bot
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target